# puppet/augeas/augeas-0.10.0/tests/root/etc/crontab
0 0 * * * root echo hello
1 0 * * * root echo hi
